# services/isen-plannings.service
[Unit]
Description=Get all plannings from ISEN Lille

[Service]
Type=simple
ExecStart=/usr/bin/GetAllPlannings.py

[Install]
WantedBy=default.target
